# templates/default/supervised-program.conf.erb
[program:<%= @program['id'] %>]
<% @program.each do |k,v| %>
<%= "#{k}=\"#{v}\"" unless k=="id" %>
<% end -%>
